# build_files/systemd/bootc-fetch.service
[Unit]
Description=Fetch bootc updates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc update --quiet
