# build_files/systemd/bootc-fetch.timer
[Unit]
Description=Fetch bootc updates weekly

[Timer]
OnCalendar=Tue *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
