# build_files/systemd/docker-groups.service
[Unit]
Description=Add wheel members to docker group

[Service]
Type=oneshot
ExecStart=/usr/bin/docker-groups

[Install]
WantedBy=default.target
