# build_files/systemd/vscode-setup.service
[Unit]
Description=VS Code first-login setup
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/vscode-setup

[Install]
WantedBy=default.target
